# contrib/udev.rules.in
ACTION!="add", GOTO="fjbtndrv_end"
SUBSYSTEM!="input", GOTO="fjbtndrv_end"

DRIVERS=="fujitsu-tablet", RUN+="@sbindir@/fjbproxy"

LABEL="fjbtndrv_end"
