# debian/postinst
#!/bin/sh
set -e

udevcontrol reload_rules
#invoke-rc.d hal restart

modprobe fsc_btns || true

